# probes/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g daemon off

[include]
files = /etc/supervisor/conf.d/*.conf
